# vms/roles/rke/templates/startup.service.j2
[Unit]
Description = startup service starting
After       = syslog.target

[Service]
ExecStart   = /root/startup.sh

[Install]
WantedBy=default.target
